# docs/supervisor.conf
# Supervisor Configuration for Kanban Task Manager
# Copy this file to /etc/supervisor/conf.d/kanban-worker.conf
# Adjust paths as needed for your installation

[program:kanban-worker]
command=/usr/bin/php /path/to/your/project/spark jobs:work --daemon --sleep=5
directory=/path/to/your/project
user=www-data
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/kanban-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=CI_ENVIRONMENT="production"

# If you want multiple workers:
# [program:kanban-worker-2]
# command=/usr/bin/php /path/to/your/project/spark jobs:work --daemon --sleep=5
# directory=/path/to/your/project
# user=www-data
# autostart=true
# autorestart=true
# redirect_stderr=true
# stdout_logfile=/var/log/supervisor/kanban-worker-2.log
# environment=CI_ENVIRONMENT="production"